# builder/workspace/files/supervisor.d/dnsmasq.ini
[program:dnsmasq]
command=dnsmasq -C /var/dnsmasq/dnsmasq.conf -k
