# builder/workspace/files/supervisor.d/openldap.ini
[program:slapd]
command=/usr/bin/slapd -u ldap -g ldap -d 1
redirect_stderr=true
