# builder/workspace/files/supervisor.d/php-fpm.ini
[program:php-fpm]
command=php-fpm -F -y /opt/php/etc/php-fpm/php-fpm.conf
